# Hw5/1651574-000107/test-1651574.service
[Unit]
Description = ShenJian homework

[Service]
Type=forking
ExecStart=/usr/sbin/test-1651574
Restart=always

[Install]
WantedBy=multi-user.target